# install/config/systemd/iqrf-gateway-webapp-mq@.service
[Unit]
Description=IQRF Gateway Webapp MQ %I
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c '/usr/sbin/iqrf-gateway-webapp-manager rabbitmq:consumer %i'
RestartSec=5
Restart=always

[Install]
WantedBy=multi-user.target